# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

exec \
        s6-setuidgid root openvpn \
        --config /etc/openvpn/vpn.conf
